# api/docker/php-fpm/supervisor/supervisor.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile          = /var/log/supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups  = 10
loglevel         = info
pidfile          = /tmp/supervisord.pid
nodaemon         = false
minfds           = 1024
minprocs         = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:ratchet]
command                 = bash -c "ulimit -n 10000; exec /usr/local/bin/php /application/bin/socket.php"
process_name            = Ratchet
numprocs                = 1
autostart               = true
autorestart             = true
user                    = root
stdout_logfile          = /var/log/supervisor/info.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /var/log/supervisor/error.log
stderr_logfile_maxbytes = 1MB